# shell/dateservice.sh
#! /bin/bash
#chkconfig: 35 66 88
#description:auto_run

# 此脚本为介绍的第二种开机自启动脚本，即直接把脚本添加至系统执行文件中
ntpdate cn.ntp.org.cn